# scripts/supervisord.conf
; ============================================================
; RadiusServer – supervisord configuration
; ============================================================
; Manages PostgreSQL, FreeRADIUS, and the Node.js backend
; inside a single container. Each process is monitored and
; restarted on failure.
; ============================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; --------------- PostgreSQL ---------------
[program:postgresql]
command=postgres -D /var/lib/postgresql/data -c listen_addresses=127.0.0.1 -c port=5432 -c shared_buffers=128MB -c max_connections=100 -c log_min_messages=warning
user=postgres
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=10
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postgresql-supervisor.log
stdout_logfile_maxbytes=5MB

; --------------- FreeRADIUS ---------------
[program:freeradius]
command=/app/scripts/freeradius-start.sh
user=radius
autostart=true
autorestart=true
startsecs=10
startretries=5
priority=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/freeradius.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
stopsignal=TERM

; --------------- Backend API ---------------
[program:backend]
command=/usr/bin/node /app/backend/index.js
user=root
autostart=true
autorestart=true
startsecs=5
startretries=3
priority=40
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10
stopsignal=TERM
environment=NODE_ENV="%(ENV_NODE_ENV)s",APP_PORT="%(ENV_APP_PORT)s",JWT_SECRET="%(ENV_JWT_SECRET)s",POSTGRES_DB="%(ENV_POSTGRES_DB)s",POSTGRES_USER="%(ENV_POSTGRES_USER)s",POSTGRES_PASSWORD="%(ENV_POSTGRES_PASSWORD)s",INTERNAL_SECRET="%(ENV_INTERNAL_SECRET)s",BCRYPT_ROUNDS="%(ENV_BCRYPT_ROUNDS)s",LOG_LEVEL="%(ENV_LOG_LEVEL)s",DEFAULT_ADMIN_USER="%(ENV_DEFAULT_ADMIN_USER)s",DEFAULT_ADMIN_PASSWORD="%(ENV_DEFAULT_ADMIN_PASSWORD)s"
